# packaging/systemd/biz-relay.service
[Unit]
Description=Biz Tunnel Relay
Wants=network-online.target
After=network-online.target

[Service]
User=biz-tunnel
Group=biz-tunnel
ExecStart=/usr/local/bin/biz-relay --config /etc/biz-tunnel/relay.toml
ExecReload=/usr/local/bin/biz-tunnelctl reload --admin http://127.0.0.1:18080
Restart=on-failure
RestartSec=3
LimitNOFILE=1048576
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/etc/biz-tunnel

[Install]
WantedBy=multi-user.target
